# coalition-worker.service
[Unit]
Description=Coalition Worker
Documentation=file:/opt/coalition-worker/README.md
After=network.target

[Service]
Type=simple
User=prod
Group=prod
WorkingDirectory=/opt/coalition-worker
ExecStart=/usr/bin/python3 /opt/coalition-worker/worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=coalition-worker

# Security settings - modified to allow access to /home/prod and /mnt
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/coalition-worker /home/prod /mnt/s3-runnable-builds

[Install]
WantedBy=multi-user.target