# coalition-worker.service
[Unit]
Description=Coalition Worker
Documentation=file:/opt/coalition-worker/README.md
After=network.target

[Service]
Type=simple
User=coalition
Group=coalition
WorkingDirectory=/opt/coalition-worker
ExecStart=/usr/bin/python3 /opt/coalition-worker/worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=coalition-worker

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/coalition-worker

[Install]
WantedBy=multi-user.target
